# infra/systemd/research-agent.service
[Unit]
Description=Daily GenAI Research Agent
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=agent
Group=agent

ExecStart=/opt/research-agent/repo/research-agent/run-agent.sh

# Secrets are encrypted at rest using a machine-specific key (via systemd-creds).
# Only systemd on this specific machine can decrypt them. Placed at deploy time.
LoadCredentialEncrypted=anthropic-api-key:/etc/research-agent/anthropic-api-key
LoadCredentialEncrypted=github-pat:/etc/research-agent/github-pat
LoadCredentialEncrypted=openalex-mailto:/etc/research-agent/openalex-mailto
LoadCredentialEncrypted=s2-api-key:/etc/research-agent/s2-api-key

# Research + critic + revision + push can take 20+ minutes
TimeoutStartSec=1800

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
