# infra/systemd/web-app.service
[Unit]
Description=Annotation web app
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=agent
Group=agent
ExecStart=/opt/research-agent/scripts/web-start.sh
ExecStop=/usr/bin/docker stop research-agent-web
Restart=on-failure
RestartSec=30

LoadCredentialEncrypted=github-pat:/etc/research-agent/github-pat

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
